# build/samba/etc/services.d/smbd/run
#!/usr/bin/with-contenv bash

nmbd -D
smbd -FS
